# files/etc/systemd/system/awg-quick@.service.d/override.conf
[Service]
ExecReload=/bin/bash -c 'exec /usr/bin/awg syncconf %i <(exec /usr/bin/awg-quick strip %i)'
