# files/etc/init.d/awg-quick
#!/sbin/openrc-run

#################################################################################
#	Filename:		files/etc/init.d/awg-quick   								#
#	Purpose:		OpenRC service that starts AmneziaWG	        			#
#	Authors:		Giulio Coa <34110430+giulioc008@users.noreply.github.com>	#
#	License:		This file is licensed under the MIT.						#
#################################################################################

description="AmneziaWG - Quick version"

depend() {
    need localmount
    need net
}

start() {
    awg-quick up wg0
}

stop() {
    awg-quick down wg0
}
